# effectors/com.vivesca.soma-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vivesca.soma-pull</string>
    <key>ProgramArguments</key>
    <array>
        <string>/home/terry/.local/share/mise/installs/python/3.13.12/bin/python3</string>
        <string>/home/terry/germline/effectors/soma-pull</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/home/terry/tmp/soma-pull-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/home/terry/tmp/soma-pull-stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
