# loci/crontab
# Vivesca crontab (supercronic on soma)
# TZ = Asia/Hong_Kong (set via supervisor.conf, matches soma-watchdog/vivesca MCP)
# Format: m h dom mon dow  command
#
# All commands SHOULD go through soma-cron-wrap for env setup + logging.
# soma-cron-wrap sources .zshenv secrets, sets PATH, captures exit codes
# to ~/.local/share/vivesca/cron-runs.jsonl.
#
# Never migrate jobs that need macOS APIs (Chrome cookies, peekaboo, iPhone
# mirroring, osascript, etc.) — those stay on Mac. See:
#   ~/docs/solutions/soma-primary-migration-plan.md

# ── Pre-existing ────────────────────────────────────────────────────
# evergreen: update Python + deps daily at 12:00 HKT
0 12 * * * soma-cron-wrap evergreen /home/vivesca/germline/effectors/evergreen /home/vivesca/germline

# supervisor-drift-check: catch /etc <-> loci/ supervisor.conf drift early.
# 28 Apr 2026 outage: tailscaled program missing in /etc, SSH-to-mac dead, no autorestart.
# Daily at 09:00 HKT — supercronic logs non-zero exit to ~/.local/share/vivesca/cron-runs.jsonl.
0 9 * * * soma-cron-wrap supervisor-drift-check /home/vivesca/germline/effectors/supervisor-drift-check

# ── Phase 2 migrations (added one at a time after end-to-end test) ──
# All added jobs were tested via `soma-cron-wrap <name> <cmd>` first.
# Mac equivalents still run in parallel until Phase 3 (DR gating) lands.

# Pilot: rotate-logs (Sunday 05:00 HKT) — low-stakes, idempotent
# Note: rotates SOMA's own logs (not Mac's) — Mac plist still rotates Mac logs.
# These are independent (different log dirs), no DR gate needed.
0 5 * * 0 soma-cron-wrap rotate-logs /usr/bin/python3 /home/vivesca/germline/effectors/rotate-logs.py

# pondus-monitor: daily model rankings check (read-only API call, idempotent)
# Mac plist also runs at 09:00 — duplicate API call until Phase 3 gating
0 9 * * * soma-cron-wrap pondus-monitor /home/vivesca/germline/.venv/bin/pondus monitor check

# wewe-rss-health: every 6h, checks WeChat RSS scraper health endpoint
0 */6 * * * soma-cron-wrap wewe-rss-health /home/vivesca/.local/bin/uv run --script --python 3.13 /home/vivesca/germline/effectors/wewe-rss-health.py

# update-coding-tools: every 2h on soma (apt + pip + uv + cargo updates)
# Mac runs hourly — soma runs less frequently (less code churn on server)
# These are independent (Mac updates Mac brew, soma updates soma apt) — no DR gate needed
0 */2 * * * soma-cron-wrap update-coding-tools /home/vivesca/germline/effectors/update-coding-tools.sh

# transfection: mirror ~/code/ git repos from soma to ganglion (every 15 min)
# Walks ~/code/*/.git, pipes owner/repo slugs to ganglion, clone-or-pulls each.
# Keeps ganglion ready as an mtor worker target for any repo Terry builds on soma.
*/15 * * * * soma-cron-wrap transfection /home/vivesca/germline/effectors/transfection

# cc-snapshot: save tmux pane → CC session ID mapping (every 2 min)
# Used by soma-cc-resume to restore sessions after restart.
*/2 * * * * python3 /home/vivesca/germline/effectors/soma-cc-snapshot

# chromatin-reader: rebuild Quartz static site from vault (every 30 min)
# Serves on port 6080 via supervisor. 9-min build time on soma — 30m interval is safe.
# DISABLED 2026-05-01 (load spike) #*/30 * * * * cd /home/vivesca/quartz && npx quartz build >> /home/vivesca/tmp/chromatin-rebuild.log 2>&1

# qmd-reindex: embed new/changed vault notes for semantic search (every 30 min)
# Self-skips via pgrep guard if a qmd embed is already running. CPU-only on soma.
*/30 * * * * /home/vivesca/germline/effectors/qmd-reindex.sh >> /home/vivesca/tmp/qmd-reindex.log 2>&1

# m3-plist-retry: deploy hardened epigenome-pull LaunchAgent to m3 when next reachable.
# Self-disables on success (~/.local/share/vivesca/m3-plist-deployed flag).
# Hard deadline 2026-05-12 — after which exit 2 surfaces via cron-runs.jsonl.
# See finding_obsidian_headless_official_cli.md and reference_terry_two_macs_imac_and_macbook.md.
0 */6 * * * soma-cron-wrap m3-plist-retry /home/vivesca/germline/effectors/m3-plist-retry
