# loci/crontab
# Vivesca crontab (supercronic)
# Times in UTC — soma runs UTC, HKT = UTC+8

# evergreen: update Python + deps daily at 04:00 UTC (12:00 HKT)
0 4 * * * /home/terry/germline/effectors/evergreen /home/terry/germline
