# loci/supervisor.conf
[supervisord]
nodaemon=true
logfile=/home/vivesca/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:vivesca]
command=/home/vivesca/germline/.venv/bin/vivesca serve --http
directory=/home/vivesca/germline
user=vivesca
autostart=true
autorestart=true
stdout_logfile=/home/vivesca/tmp/mcp-server.log
stderr_logfile=/home/vivesca/tmp/mcp-server.log
environment=HOME="/home/vivesca",TZ="Asia/Hong_Kong",PATH="/home/vivesca/germline/.venv/bin:/home/vivesca/germline/effectors:/home/vivesca/.local/bin:/usr/local/bin:/usr/bin:/bin"

[program:soma-watchdog]
command=python3 /home/vivesca/germline/effectors/soma-watchdog
directory=/home/vivesca/germline
user=vivesca
autostart=true
autorestart=true
stdout_logfile=/home/vivesca/tmp/soma-watchdog.log
stderr_logfile=/home/vivesca/tmp/soma-watchdog.log
environment=HOME="/home/vivesca",TZ="Asia/Hong_Kong",PATH="/home/vivesca/germline/.venv/bin:/home/vivesca/germline/effectors:/home/vivesca/.local/bin:/usr/local/bin:/usr/bin:/bin"

[program:soma-liveness]
command=python3 /home/vivesca/germline/effectors/soma-liveness
directory=/home/vivesca
user=vivesca
autostart=true
autorestart=true
stdout_logfile=/home/vivesca/tmp/soma-liveness.log
stderr_logfile=/home/vivesca/tmp/soma-liveness.log
environment=HOME="/home/vivesca",TZ="Asia/Hong_Kong"

[program:supercronic]
command=/home/vivesca/.local/bin/supercronic -passthrough-logs /home/vivesca/germline/loci/crontab
directory=/home/vivesca
user=vivesca
autostart=true
autorestart=true
stdout_logfile=/home/vivesca/tmp/supercronic.log
stderr_logfile=/home/vivesca/tmp/supercronic.log
environment=HOME="/home/vivesca",TZ="Asia/Hong_Kong",PATH="/home/vivesca/germline/effectors:/home/vivesca/germline/.venv/bin:/home/vivesca/.local/bin:/home/vivesca/.cargo/bin:/home/vivesca/.bun/bin:/usr/local/bin:/usr/bin:/bin"

[program:ob-sync]
command=/home/vivesca/.local/bin/ob sync --continuous --path /home/vivesca/epigenome/chromatin
directory=/home/vivesca
user=vivesca
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/home/vivesca/tmp/ob-sync.log
stderr_logfile=/home/vivesca/tmp/ob-sync.log
environment=HOME="/home/vivesca",TZ="Asia/Hong_Kong",PATH="/home/vivesca/.local/bin:/usr/local/bin:/usr/bin:/bin"

[program:temporal-tunnel]
command=ssh -N -L 7233:localhost:7233 -o ServerAliveInterval=60 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes ganglion
directory=/home/vivesca
user=vivesca
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/home/vivesca/tmp/temporal-tunnel.log
stderr_logfile=/home/vivesca/tmp/temporal-tunnel.log
environment=HOME="/home/vivesca"

[program:tailscaled]
command=tailscaled --tun=userspace-networking --socks5-server=localhost:1055 --outbound-http-proxy-listen=localhost:1056 --state=/home/vivesca/.tailscale/tailscaled.state
directory=/
user=root
autostart=true
autorestart=true
stdout_logfile=/home/vivesca/tmp/tailscaled.log
stderr_logfile=/home/vivesca/tmp/tailscaled.log
