# loci/supervisor.conf
[supervisord]
nodaemon=true
logfile=/home/terry/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:vivesca]
; TODO: after next `fly deploy` or machine restart, remove supervisor-env wrapper
; (Fly secrets now staged for EXA/TAVILY/SERPER/JINA/FIRECRAWL/XAI/PERPLEXITY)
command=/home/terry/germline/effectors/supervisor-env /home/terry/germline/.venv/bin/vivesca serve --http
directory=/home/terry/germline
user=terry
autostart=true
autorestart=true
stdout_logfile=/home/terry/tmp/mcp-server.log
stderr_logfile=/home/terry/tmp/mcp-server.log

[program:soma-watchdog]
command=python3 /home/terry/germline/effectors/soma-watchdog
directory=/home/terry/germline
user=terry
autostart=true
autorestart=true
stdout_logfile=/home/terry/tmp/soma-watchdog.log
stderr_logfile=/home/terry/tmp/soma-watchdog.log
environment=HOME="/home/terry",TZ="Asia/Hong_Kong",PATH="/home/terry/germline/.venv/bin:/home/terry/.local/bin:/usr/local/bin:/usr/bin:/bin"

[program:supercronic]
command=/home/terry/.local/bin/supercronic /home/terry/germline/loci/crontab
directory=/home/terry/germline
user=terry
autostart=true
autorestart=true
stdout_logfile=/home/terry/tmp/supercronic.log
stderr_logfile=/home/terry/tmp/supercronic.log
environment=HOME="/home/terry",TZ="UTC",PATH="/home/terry/germline/.venv/bin:/home/terry/germline/effectors:/home/terry/.local/bin:/home/terry/.cargo/bin:/usr/local/bin:/usr/bin:/bin"

[program:tailscaled]
command=tailscaled --state=/tmp/tailscale.state
directory=/
user=root
autostart=true
autorestart=true
stdout_logfile=/home/terry/tmp/tailscaled.log
stderr_logfile=/home/terry/tmp/tailscaled.log
